# drink-water-reminder.service
[Unit]
Description=Hydration reminder
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/drink-water-reminder
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target